# os/rootfs/usr/lib/systemd/system/novanas-overlay-etc.mount
[Unit]
Description=NovaNas overlay for /etc
DefaultDependencies=no
Conflicts=umount.target
After=novanas-persistent.mount
Requires=novanas-persistent.mount
Before=local-fs.target sysinit.target

[Mount]
What=overlay
Where=/etc
Type=overlay
# lowerdir is the read-only /etc shipped in the image; upper+work come from
# the persistent partition. ExecStartPre creates the upper/work dirs if this
# is a fresh install.
Options=lowerdir=/etc,upperdir=/mnt/persistent/etc,workdir=/mnt/persistent/etc.work

[Install]
WantedBy=local-fs.target
