# os/rootfs/etc/systemd/system/k3s.service.d/override.conf
# NovaNas-specific k3s tuning. The base k3s.service is installed by the k3s
# installer; we override knobs here so updates to the upstream unit don't
# clobber our config.

[Service]
# Data under the persistent overlay so upgrades preserve cluster state.
Environment=K3S_DATA_DIR=/var/lib/rancher/k3s
Environment=K3S_KUBECONFIG_MODE=0600

# Disable components we replace or don't use:
#   servicelb  - replaced by novaedge
#   traefik    - we ship our own ingress via novaedge
#   local-storage - we use the NovaNas CSI
#   metrics-server - shipped via LGTM subchart
ExecStart=
ExecStart=/usr/local/bin/k3s server \
  --disable=servicelb \
  --disable=traefik \
  --disable=local-storage \
  --disable=metrics-server \
  --kubelet-arg=max-pods=250 \
  --node-label=novanas.io/role=control-plane \
  --write-kubeconfig-mode=0600 \
  --embedded-registry \
  --tls-san=nas.local

Restart=on-failure
RestartSec=5s
# k3s is the critical path for everything; wait long enough that a slow disk
# doesn't trigger rollback.
TimeoutStartSec=300s
