# deploy/systemd/rpc-gssd.service.d/override.conf
# /etc/systemd/system/rpc-gssd.service.d/override.conf
#
# Drop-in for nfs-common's rpc-gssd.service. Two changes vs. the
# upstream Debian unit:
#
#  1. -n: rpc.gssd uses the host's /etc/krb5.keytab (specifically the
#     nfs/<host>@REALM service principal) for upcalls when no user TGT
#     is available. This is the "machine credentials" mode that lets
#     containerized workloads mount sec=krb5p shares without each pod
#     carrying its own TGT.
#
#  2. -f: foreground. Required so systemd's Type=notify accounting sees
#     rpc.gssd's lifecycle correctly (the upstream unit type is
#     forking; we run it in the foreground for clean restarts and
#     journal capture).
#
# The host MUST have a populated /etc/krb5.keytab containing
# nfs/<fqdn>@<REALM> before this unit can start. NovaNAS's
# bootstrap script (deploy/krb5/nova-kdc-bootstrap.sh, owned by
# Agent A) provisions the keytab.

[Service]
ExecStart=
ExecStart=/usr/sbin/rpc.gssd -n -f
Restart=on-failure
RestartSec=5s
