# os/rootfs/etc/systemd/system/novanas-installer-watchdog.service
[Unit]
Description=NovaNas live-boot poweroff watchdog
# Belt-and-suspenders: even if novanas-installer.service never runs or
# hangs, this watchdog powers the VM off after 5 minutes so packer's
# shutdown_timeout catches a clean shutdown rather than timing out.
# Happy-path dry-run install completes in <1 minute, so 5 minutes is
# plenty of slack while keeping CI feedback tight.
# Only active on live boots that carry our installer flag — never on
# installed systems.
ConditionKernelCommandLine=novanas.installer
ConditionKernelCommandLine=boot=live
# Run earlier than multi-user.target so we don't depend on the full
# boot succeeding — if the live-init hangs before multi-user.target,
# this unit still has to fire. basic.target is reached very early.
After=basic.target
DefaultDependencies=no

[Service]
Type=oneshot
# DMI sys_vendor signals VM vs bare metal. On real hardware the operator
# is interactively driving the installer TUI and may legitimately take
# 30+ minutes; auto-powering-off would yank the rug. Only the CI/VM
# runs need the timer (Packer's shutdown_timeout depends on it).
ExecStart=/bin/bash -c ' \
  vendor=$(cat /sys/class/dmi/id/sys_vendor 2>/dev/null || echo unknown); \
  case "$vendor" in \
    QEMU|"Red Hat"|Bochs|innotek*|"VMware, Inc."|"Microsoft Corporation") \
      echo "novanas-installer-watchdog: VM detected ($vendor) — sleeping 5m then poweroff" | tee /dev/console; \
      sleep 300; \
      echo "novanas-installer-watchdog: timeout — forcing poweroff" | tee /dev/console; \
      /sbin/poweroff -f || (echo o > /proc/sysrq-trigger) ;; \
    *) \
      echo "novanas-installer-watchdog: real hardware ($vendor) — auto-poweroff disabled" | tee /dev/console ;; \
  esac'

[Install]
WantedBy=basic.target
