# os/rootfs/etc/systemd/system/novanas-installer.service
[Unit]
Description=NovaNas unattended installer (live boot)
Documentation=https://github.com/azrtydxb/novanas
# Only run when booted from the installer ISO. live-boot sets boot=live
# on the kernel cmdline; we also require our own signal so a rescue boot
# doesn't accidentally trigger an install.
ConditionKernelCommandLine=novanas.installer
ConditionKernelCommandLine=boot=live
After=systemd-udev-settle.service local-fs.target
Wants=systemd-udev-settle.service
# Block any getty from stealing focus while we're running.
Conflicts=getty@tty1.service
After=getty@tty1.service

[Service]
Type=oneshot
RemainAfterExit=no
# Route output to both journal and the serial console so packer VNC /
# virsh console surface installer progress.
StandardOutput=journal+console
StandardError=journal+console
TTYPath=/dev/tty1
ExecStart=/usr/local/bin/novanas-live-installer.sh
TimeoutStartSec=25min

[Install]
WantedBy=multi-user.target
