# deploy/systemd/keycloak.service.d/management-port.conf
# Move Keycloak's Quarkus management endpoint off the default port 9000
# because RustFS object storage (deploy/systemd/rustfs.service) claims it
# as the conventional MinIO/S3-compatible API port. The management
# endpoint exposes /q/health/, /q/metrics, and other non-public surfaces;
# 9095 is unused elsewhere in the NovaNAS port allocation.
#
# Verified on Keycloak 26.0.7. KC_HTTP_MANAGEMENT_PORT was introduced in
# Keycloak 25; older releases don't separate management from public.
[Service]
Environment=KC_HTTP_MANAGEMENT_PORT=9095
