# installer/assets/systemd/mnt-persistent.mount
[Unit]
Description=NovaNas persistent partition
DefaultDependencies=no
Conflicts=umount.target
Before=local-fs.target umount.target
# Deliberately NOT listing systemd-fsck@dev-disk-by\x2dpartlabel-Persistent
# in After=: on a live boot that ordering would cause systemd to
# materialise the device unit and wait 90s for it to appear before the
# mount's Condition is even evaluated. fsck still runs at mount time
# via the generator's default path.
# Skip on live boot — the Persistent partition only exists on a
# disk-installed system.
ConditionKernelCommandLine=!boot=live
ConditionPathExists=/dev/disk/by-partlabel/Persistent

[Mount]
What=/dev/disk/by-partlabel/Persistent
Where=/mnt/persistent
Type=ext4
# nofail: don't fail local-fs.target if the partition is missing (live boot).
# x-systemd.device-timeout=1: give up on the device unit in 1s instead of
# the default 90s — critical on live boot where the partition never
# exists.
Options=rw,noatime,errors=remount-ro,nofail,x-systemd.device-timeout=1s

[Install]
WantedBy=local-fs.target
